# ops/supervisor/conf.d/uwsgi.conf
[program:uwsgi]
command=uwsgi --ini uwsgi.ini
autostart=true
autorestart=true
startsecs=10
startretries=3
stopsignal=QUIT
user=root
stdout_logfile = /tmp/uwsgi.stdout.log
stderr_logfile = /tmp/uwsgi.stderr.log
directory=/usr/src/app
environment=DJANGO_SETTINGS_MODULE=config.prod